# root/etc/s6-overlay/s6-rc.d/update-stats/run
#!/usr/bin/bash
# Service to update statistics cache every 15 minutes

# Wait for calibre-web to be ready
sleep 60

while true; do
    echo "$(date) - Updating statistics cache..."
    cd /app/calibre-web && /lsiopy/bin/python -c "from cps import admin_stats; admin_stats.update_stats_cache()" || echo "Failed to update stats cache"

    # Wait 15 minutes
    sleep 900
done
